# build/package/files/etc/init.d/api-umbrella
#!/usr/bin/env bash
#
# api-umbrella
#
# chkconfig: - 85 15
# description: API Umbrella
# processname: api-umbrella
# config: /etc/api-umbrella/api-umbrella.yml

### BEGIN INIT INFO
# Provides:          api-umbrella
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: API Umbrella
### END INIT INFO

name=api-umbrella
desc="API Umbrella"

if [ -f /etc/rc.d/init.d/functions ]; then
  # shellcheck disable=SC1091
  . /etc/rc.d/init.d/functions
elif [ -f /lib/lsb/init-functions ]; then
  # shellcheck disable=SC1091
  . /lib/lsb/init-functions
fi

if [ -f /etc/sysconfig/api-umbrella ]; then
  # shellcheck disable=SC1091
  . /etc/sysconfig/api-umbrella
elif [ -f /etc/default/api-umbrella ]; then
  # shellcheck disable=SC1091
  . /etc/default/api-umbrella
fi

start() {
  if type log_daemon_msg > /dev/null 2>&1; then
    log_daemon_msg "Starting $desc" "$name"
  else
    echo -n $"Starting $name: "
  fi

  api-umbrella start
  retval=$?

  if type log_end_msg > /dev/null 2>&1; then
    log_end_msg $retval
  elif type success > /dev/null 2>&1; then
    if [ $retval -eq 0 ]; then
      success $"$name startup"
    else
      failure $"$name startup"
    fi
    echo
  fi

  return $retval
}

stop() {
  if type log_daemon_msg > /dev/null 2>&1; then
    log_daemon_msg "Stopping $desc" "$name"
  else
    echo -n $"Stopping $name: "
  fi

  api-umbrella stop
  retval=$?

  if type log_end_msg > /dev/null 2>&1; then
    log_end_msg $retval
  elif type success > /dev/null 2>&1; then
    if [ $retval -eq 0 ]; then
      success $"$name shutdown"
    else
      failure $"$name shutdown"
    fi
    echo
  fi

  return $retval
}

restart() {
  stop
  start
}

reload() {
  if type log_daemon_msg > /dev/null 2>&1; then
    log_daemon_msg "Reloading $desc" "$name"
  else
    echo -n $"Reloading $name: "
  fi

  api-umbrella reload
  retval=$?

  if type log_end_msg > /dev/null 2>&1; then
    log_end_msg $retval
  elif type success > /dev/null 2>&1; then
    if [ $retval -eq 0 ]; then
      success $"$name reload"
    else
      failure $"$name reload"
    fi
    echo
  fi

  return $retval
}

status() {
  api-umbrella status
  retval=$?
  return $retval
}

status_quiet() {
  status > /dev/null 2>&1
}

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  status)
    status
    ;;
  restart)
    restart
    ;;
  reload)
    reload
    ;;
  condrestart)
    status_quiet || exit 0
    restart
    ;;
  *)
    echo "Usage: $name {start|stop|status|reload|restart|condrestart}"
    exit 1
    ;;
esac
exit $?
