# build/package/scripts/after-remove
#!/usr/bin/env bash

set -e -u

remove_accounts=false
remove_data=false
remove_files=false
remove_service=false
restart_service=false

case "$1" in
  # dpkg
  remove)
    remove_files=true
    remove_service=true
    ;;
  purge)
    remove_accounts=true
    remove_data=true
    ;;
  abort-install|abort-upgrade|disappear|failed-upgrade|upgrade)
    ;;

  # rpm
  0) # uninstall
    remove_files=true
    ;;
  1) # upgrade
    restart_service=true
    ;;
esac

if [ "$remove_files" = "true" ]; then
  rm -f /usr/bin/api-umbrella /var/log/api-umbrella
  rmdir /etc/api-umbrella > /dev/null 2>&1 || true # Delete directory if empty.
fi

if [ "$remove_service" = "true" ]; then
  if command -v update-rc.d > /dev/null 2>&1; then
    update-rc.d api-umbrella remove > /dev/null
  fi
fi

if [ "$remove_data" = "true" ]; then
  rm -rf /etc/api-umbrella /opt/api-umbrella
fi

if [ "$remove_accounts" = "true" ]; then
  user=api-umbrella
  group=api-umbrella
  deploy_user=api-umbrella-deploy
  deploy_group=api-umbrella-deploy

  if getent passwd $deploy_user > /dev/null; then
    userdel $deploy_user
  fi

  if getent passwd $user > /dev/null; then
    userdel $user
  fi

  if getent group $deploy_group > /dev/null; then
    groupdel $deploy_group
  fi

  if getent group $group > /dev/null; then
    groupdel $group
  fi
fi

if [ "$restart_service" = "true" ]; then
  # On upgrades, restart the service if it's currently running.
  /etc/init.d/api-umbrella condrestart
fi
