# build/package/scripts/after-remove
#!/usr/bin/env bash

set -e -u

remove_accounts=false
remove_data=false
remove_files=false
remove_service=false
restart_service=false

case "$1" in
  # dpkg
  remove)
    remove_files=true
    remove_service=true
    ;;
  purge)
    remove_accounts=true
    remove_data=true
    ;;
  abort-install|abort-upgrade|disappear|failed-upgrade|upgrade)
    ;;

  # rpm
  0) # uninstall
    remove_files=true
    ;;
  1) # upgrade
    restart_service=true
    ;;
esac

destdir="${DESTDIR:-}"
prefix_dir="$destdir${PREFIX:-/opt/api-umbrella}"

if [ "$remove_files" = "true" ]; then
  rm -f "$destdir/usr/bin/api-umbrella" "$destdir/var/log/api-umbrella"
  rmdir "$destdir/etc/api-umbrella" > /dev/null 2>&1 || true # Delete directory if empty.
fi

if [ "$remove_service" = "true" ]; then
  if command -v update-rc.d > /dev/null 2>&1; then
    update-rc.d api-umbrella remove > /dev/null
  fi
fi

if [ "$remove_data" = "true" ]; then
  rm -rf "$destdir/etc/api-umbrella" "$prefix_dir"
fi

if [ "$remove_accounts" = "true" ]; then
  user=api-umbrella
  group=api-umbrella
  auto_ssl_user=api-umbrella-auto-ssl
  auto_ssl_group=api-umbrella-auto-ssl

  if getent passwd $auto_ssl_user > /dev/null; then
    userdel $auto_ssl_user
  fi

  if getent passwd $user > /dev/null; then
    userdel $user
  fi

  if getent group $auto_ssl_group > /dev/null; then
    groupdel $auto_ssl_group
  fi

  if getent group $group > /dev/null; then
    groupdel $group
  fi
fi

if [ "$restart_service" = "true" ]; then
  # On upgrades, restart the service if it's currently running.
  "$destdir/etc/init.d/api-umbrella" condrestart
fi
