# build/package/scripts/after-remove
#!/bin/bash
#
# Perform necessary api-umbrella removal steps
# after package is uninstalled.
#

function error_exit {
  echo "api-umbrella: ${1:-"Unknown Error"}" 1>&2
  exit 1
}

function remove_service {
  if command -v chkconfig > /dev/null 2>&1; then
    chkconfig --del api-umbrella > /dev/null 2>&1
  elif command -v update-rc.d > /dev/null 2>&1; then
    update-rc.d api-umbrella remove > /dev/null 2>&1
  fi

  rm -f /usr/bin/api-umbrella
  rm -f /var/log/api-umbrella
  rmdir /etc/api-umbrella > /dev/null 2>&1 # Delete directory if empty.
}

# On RedHat RPM systems, when an upgrade is performed, the old package is
# removed after the newer one is installed. Only perform cleanup on full
# removals.
#
# See: http://tickets.opscode.com/browse/CHEF-3022
if [ ! -f /etc/redhat-release -a ! -f /etc/fedora-release -a ! -f /etc/system-release ]; then
  # Non-RPM system: Go ahead with removal
  remove_service
elif [ "$1" -eq "0" ]; then
  # RPM system: Uninstalling (not upgrading)
  remove_service
elif [ "$1" -eq "1" ]; then
  # RPM system: Upgrading (not uninstalling)
  # Restart, but only if started.
  /etc/init.d/api-umbrella status > /dev/null 2>&1 || exit 0
  /etc/init.d/api-umbrella restart > /dev/null 2>&1 || error_exit "Cannot restart api-umbrella service"
fi

exit 0
