# build/package/scripts/before-remove
#!/bin/bash

set -e -u

STOP_SERVICE=false
REMOVE_SERVICE=false

case "$1" in
  # dpkg
  remove)
    STOP_SERVICE=true
    ;;
  deconfigure|failed-upgrade|upgrade)
    ;;

  # rpm
  0) # uninstall
    STOP_SERVICE=true
    REMOVE_SERVICE=true
    ;;
  1) # upgrade
    ;;
esac

if [ "$STOP_SERVICE" = "true" ]; then
  /etc/init.d/api-umbrella stop || true
fi

if [ "$REMOVE_SERVICE" = "true" ]; then
  if command -v chkconfig > /dev/null 2>&1; then
    chkconfig --del api-umbrella
  fi
fi
